# nauczsiepythona.cron.d
# renew ssl cert
0 5 1 */2 * /usr/local/bin/docker-compose -f /home/tryb/Aplikacja-do-nauki-Pythona/docker-compose.yml up -d certbot
